# deploy/collector/launchd/com.tokemetry.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd agent for the tokemetry collector (macOS).
     Install to ~/Library/LaunchAgents/ and edit the paths for your user.
     Load with: launchctl load ~/Library/LaunchAgents/com.tokemetry.collector.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tokemetry.collector</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/CHANGE_ME/.local/bin/tokemetry-collector</string>
      <string>--config</string>
      <string>/Users/CHANGE_ME/.config/tokemetry/collector.toml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/CHANGE_ME/Library/Logs/tokemetry-collector.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/CHANGE_ME/Library/Logs/tokemetry-collector.log</string>
  </dict>
</plist>
